# stats.service.example
[Unit]
Description=Shalog stats service

[Service]
WorkingDirectory=/home/tracking/shalog/
ExecStart=/home/tracking/shalog/stats.py

[Install]
WantedBy=default.target
